# systemd/pi-voice-assistant.service
[Unit]
Description=pi-voice-assistant wake word daemon ("Alexa" -> record question -> Claude -> spoken reply)
# No explicit audio dependency needed -- PipeWire starts on demand via socket
# activation when this (or any client) first connects to it.

[Service]
Type=simple
WorkingDirectory=/home/<PI_USER>/<PI_DIR>
ExecStart=/home/<PI_USER>/.local/bin/uv run wake_word_daemon.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
